# alpaca-bot.service
[Unit]
Description=Alpaca Trading Bot
Documentation=https://github.com/grounzero/alpaca-fleece#readme
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Use %h for home directory (systemd specifier)
WorkingDirectory=%h/.openclaw/workspace/alpaca-fleece
Environment=PYTHONUNBUFFERED=1
Environment=BOT_ROOT=%h/.openclaw/workspace/alpaca-fleece
Environment=PATH=%h/.openclaw/workspace/alpaca-fleece/.venv/bin:/usr/local/bin:/usr/bin:/bin

# Use the virtual environment Python
ExecStart=%h/.openclaw/workspace/alpaca-fleece/.venv/bin/python -u orchestrator.py

# Graceful shutdown on stop
TimeoutStopSec=30
KillSignal=SIGTERM

# Restart policy - don't restart on clean exit (0), but restart on crashes
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:%h/.openclaw/workspace/alpaca-fleece/logs/orchestrator.out
StandardError=append:%h/.openclaw/workspace/alpaca-fleece/logs/orchestrator.out

# Security hardening (optional but recommended)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=%h/.openclaw/workspace/alpaca-fleece/data %h/.openclaw/workspace/alpaca-fleece/logs

[Install]
WantedBy=default.target
