# alpaca-bot.service
[Unit]
Description=Alpaca Trading Bot
Documentation=https://github.com/grounzero/alpaca-fleece#readme
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Use %h for home directory (systemd specifier)
# Allow administrators to override environment via /etc/default/alpaca-bot
EnvironmentFile=-/etc/default/alpaca-bot

# Set a default BOT_ROOT relative to the service user's home directory
Environment=BOT_ROOT=%h/.openclaw/workspace/alpaca-fleece
Environment=PYTHONUNBUFFERED=1

# Ensure PATH includes the virtualenv inside BOT_ROOT
Environment=PATH=${BOT_ROOT}/.venv/bin:/usr/local/bin:/usr/bin:/bin

# Use the virtual environment Python inside BOT_ROOT
ExecStart=${BOT_ROOT}/.venv/bin/python -u ${BOT_ROOT}/orchestrator.py

# Graceful shutdown on stop
TimeoutStopSec=30
KillSignal=SIGTERM

# Restart policy - don't restart on clean exit (0), but restart on crashes
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:${BOT_ROOT}/logs/orchestrator.out
StandardError=append:${BOT_ROOT}/logs/orchestrator.out

# Security hardening (optional but recommended)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=${BOT_ROOT}/data ${BOT_ROOT}/logs

[Install]
WantedBy=default.target
