# alpaca-fleece/alpaca-bot.service
[Unit]
Description=Alpaca Trading Bot
Documentation=https://github.com/alpaca-trade-api
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/t-rox/.openclaw/workspace/alpaca-fleece
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/t-rox/.openclaw/workspace/alpaca-fleece/.venv/bin:/usr/local/bin:/usr/bin:/bin

# Use the virtual environment Python
ExecStart=/home/t-rox/.openclaw/workspace/alpaca-fleece/.venv/bin/python -u orchestrator.py

# Graceful shutdown on stop
TimeoutStopSec=30
KillSignal=SIGTERM

# Restart policy - don't restart on clean exit (0), but restart on crashes
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/home/t-rox/.openclaw/workspace/alpaca-fleece/logs/orchestrator.out
StandardError=append:/home/t-rox/.openclaw/workspace/alpaca-fleece/logs/orchestrator.out

# Security hardening (optional but recommended)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/t-rox/.openclaw/workspace/alpaca-fleece/data /home/t-rox/.openclaw/workspace/alpaca-fleece/logs

[Install]
WantedBy=default.target
